# supervisor/prod.conf
[program:webserver]
directory = /opt/cms
command = gunicorn project.wsgi --bind=0.0.0.0:80  --workers=5
